# base/conf/cuda/shell/etc/profile.d/00-reset-path.sh
# Reset default PATH settings
PATH="/opt/code-server/bin:/usr/local/nvidia/bin:/usr/local/cuda/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
export PATH
